# dot_config/systemd/user/fish-plug-update.service
[Unit]
Description=Update fish plugins
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/fish -c 'fisher update'
